# scripts/hermes/launchd/pro/co.jovie.hermes.cron-agentcookie-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  agentcookie-sync launchd plist — runs on the SOURCE Mac (daily-driver MBP).

  Security gate: does nothing without AGENTCOOKIE_LIVE=1 in the environment.
  Bootstrap: rendered by scripts/hermes/bootstrap-pro-launchd.sh.
  Tim approval required before loading: launchctl load ~/Library/LaunchAgents/co.jovie.hermes.cron-agentcookie-sync.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.jovie.hermes.cron-agentcookie-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{TSX_BIN}}</string>
        <string>{{JOVIE_REPO}}/scripts/hermes/jobs/agentcookie-sync.ts</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_HOME</key>
        <string>{{HOME}}/.hermes</string>
        <key>AGENTCOOKIE_TARGET_HOST</key>
        <string>hermes-air</string>
        <key>PATH</key>
        <string>{{HOME}}/.bun/bin:{{HOME}}/.hermes/bin:{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/cron-agentcookie-sync.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/cron-agentcookie-sync.err.log</string>
</dict>
</plist>
