# scripts/hermes/launchd/pro/co.jovie.hermes.agentcookie-sender.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.jovie.hermes.agentcookie-sender</string>
    <!--
        agentcookie send daemon — Pro (daily-driver Mac) side.

        Continuously watches Chrome cookie store and syncs allowed
        cookies/API keys to the Air over Tailscale (AES-256-GCM encrypted).

        Requires env vars (set in Doppler jovie-web/dev):
          AGENTCOOKIE_AIR_IP       — Tailscale IP of the Air (run: tailscale ip -4
                                     on the Air to get this)
          AGENTCOOKIE_ENCRYPT_KEY  — AES-256-GCM shared key (same key as Air)
          AGENTCOOKIE_PORT         — optional, defaults to 17070
    -->
    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.hermes/bin/agentcookie</string>
        <string>send</string>
        <string>--host-env</string>
        <string>AGENTCOOKIE_AIR_IP</string>
        <string>--port</string>
        <string>17070</string>
        <string>--encrypt-key-env</string>
        <string>AGENTCOOKIE_ENCRYPT_KEY</string>
        <string>--allowlist</string>
        <string>{{HOME}}/.hermes/state/agentcookie/allowlist.txt</string>
        <string>--blocklist</string>
        <string>{{HOME}}/.hermes/state/agentcookie/blocklist.txt</string>
        <string>--continuous</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_HOME</key>
        <string>{{HOME}}/.hermes</string>
        <key>PATH</key>
        <string>{{HOME}}/.grok/bin:{{NODE_BIN_DIR}}:{{HOME}}/.bun/bin:{{HOME}}/.hermes/bin:{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/agentcookie-sender.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/agentcookie-sender.err.log</string>
    <!-- Throttle restart if binary or env vars are missing. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
