# .github/runner-host/ci-runner-capacity-reconcile.timer
[Unit]
Description=Continuously enforce Jovie CI runner process capacity

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=ci-runner-capacity-reconcile.service

[Install]
WantedBy=timers.target
