# scripts/worktree-lifecycle-launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.jovie.worktree-lifecycle</string>
  <key>ProgramArguments</key><array>
    <string>/Users/timwhite/Jovie/scripts/worktree-lifecycle-cron.sh</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/jovie-worktree-lifecycle.out</string>
  <key>StandardErrorPath</key><string>/tmp/jovie-worktree-lifecycle.err</string>
</dict></plist>
