# scripts/hermes/launchd/co.jovie.hermes.agentcookie-receiver.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.jovie.hermes.agentcookie-receiver</string>
    <!--
        agentcookie receive daemon — Air side.

        Listens on Tailscale-only port {{AGENTCOOKIE_PORT}} for encrypted
        cookie/credential payloads from the Pro (daily-driver Mac).

        Requires env vars (set in Doppler jovie-web/dev):
          AGENTCOOKIE_ENCRYPT_KEY  — AES-256-GCM shared key (generate with:
                                     openssl rand -base64 32)
          AGENTCOOKIE_PORT         — optional, defaults to 17070

        Security: --tailscale-only flag rejects connections from non-Tailscale
        source IPs. The port is NOT opened on public interfaces.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.hermes/bin/agentcookie</string>
        <string>receive</string>
        <string>--port</string>
        <string>17070</string>
        <string>--decrypt-key-env</string>
        <string>AGENTCOOKIE_ENCRYPT_KEY</string>
        <string>--output-dir</string>
        <string>{{HOME}}/.hermes/cookies</string>
        <string>--allowlist</string>
        <string>{{HOME}}/.hermes/state/agentcookie/allowlist.txt</string>
        <string>--blocklist</string>
        <string>{{HOME}}/.hermes/state/agentcookie/blocklist.txt</string>
        <string>--tailscale-only</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_HOME</key>
        <string>{{HOME}}/.hermes</string>
        <key>PATH</key>
        <string>{{HOME}}/.bun/bin:{{HOME}}/.hermes/bin:{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/agentcookie-receiver.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.hermes/logs/launchd/agentcookie-receiver.err.log</string>
    <!-- Throttle restart if agentcookie binary is missing or misconfigured. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
